# deploy/nexus-agent.service
[Unit]
Description=Nexus Agent — per-machine daemon for Claude Code session tracking
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/nexus-agent
Restart=always
RestartSec=5
Environment=RUST_LOG=info
Environment=PATH=%h/.local/bin:%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
